# templates/redis_init_script.erb
#!/bin/sh
# chkconfig: 345 85 60
# description: Redis server - start, stop 
### BEGIN INIT INFO
# Required-Start: $local_fs $network $syslog
# Default-Start: 3 4 5
# Default-Stop: 0 1 2 6
# Short-Description: Redis server
### END INIT INFO
#
# HEADER: This file is deployed from Puppet module redis
# HEADER: Manual changes to this file will be overwritten 
# 

REDISHOST=<%= bind %>
REDISPORT=<%= port %>
EXEC=<%= deploymentdir %>/bin/redis-server
CLIEXEC=<%= deploymentdir %>/bin/redis-cli

PIDFILE=/var/run/<%= name %>.pid
CONF=<%= config_file_path %> 

case "$1" in
    start)
        if [ -f $PIDFILE ]
        then
                echo "$PIDFILE exists, process is already running or crashed"
        else
                echo "Starting Redis server..."
                $EXEC $CONF
        fi
        ;;
    stop)
        if [ ! -f $PIDFILE ]
        then
                echo "$PIDFILE does not exist, process is not running"
        else
                PID=$(cat $PIDFILE)
                echo "Stopping ..."
                $CLIEXEC -p $REDISPORT -h $REDISHOST shutdown
                while [ -x /proc/${PID} ]
                do
                    echo "Waiting for Redis to shutdown ..."
                    sleep 1
                done
                echo "Redis stopped"
        fi
        ;;
    *)
        echo "Please use start or stop as first argument"
        ;;
esac
